# debian/openssh-server.conffiles
/etc/init.d/ssh
/etc/pam.d/ssh
/etc/default/ssh
